# debian/postinst
#!/bin/sh

set -e


case "$1" in
    configure)
    
 	rm -rf /etc/zero-server-wizard/enabled.d/* || true
	rm -rf /etc/zero-server-wizard/disabled.d/* || true
	cp /usr/share/zero-server-wizard/confs/* /etc/zero-server-wizard/enabled.d/ || true   
	
	if dpkg --compare-versions "$2" lt 0.9.10; then
		CNAME="/var/lib/dnsmasq/config/cname-server"
			
		if [ -f $CNAME ]; then
		
			DOMAIN=$(n4d-vars getvalues INTERNAL_DOMAIN | cut -d"'" -f2)
			SRVNAME=$(n4d-vars getvalues HOSTNAME | cut -d"'" -f2)
			echo "Fixing ipxeboot DNS ..."
			sed -i -e "/^cname=ipxeboot.*/d" $CNAME
			sed -i -e "/^$/d" $CNAME
			echo "cname=ipxeboot.$DOMAIN,$SRVNAME.$DOMAIN\n" >> $CNAME
			invoke-rc.d dnsmasq force-reload || true
			
			zero-center set-configured zero-server-wizard || true
			
		fi
		
	fi	

	if dpkg --compare-versions "$2" lt 0.9.20; then
	
		SRV_IP="/var/lib/n4d/variables-dir/SRV_IP"
	
		if [ -f $SRV_IP ]; then
	
			/usr/share/zero-server-wizard/bin/appendsearch.py || true
			
		fi
	fi
	deb-systemd-invoke restart n4d || true
	
    
    ;;



    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
